# builder/local/crontab
# Jalankan cronjob setiap menit (di dalam container path: /var/www/html)
# Gunakan shell eksplisit agar tidak bergantung pada bit executable di cronjob.sh
* * * * * root /bin/bash /var/www/html/cronjob.sh
